# services/runit-usr/run
#!/bin/sh
#
# Usage:
# `ln -s /etc/sv/runit-usr /var/service/[someuser]`
#

UID=`pwd -P`
UID=${UID##*/}

if [ -d "/home/$UID/service" ]; then
	chpst -u$UID runsvdir /home/$UID/service
fi
