# services/runit-usr/run
#!/bin/sh

UID=`pwd -P`
UID=${UID##*/}

if [ -d "/home/$UID/service" ]; then
	chpst -u$UID runsvdir /home/$UID/service
fi
